# docs/com.contentmachine.autobroll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contentmachine.autobroll</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tarungupta/miniconda3/envs/content_engine_env/bin/python3</string>
        <string>/Users/tarungupta/Making It Big/Claude/content-machine/scripts/auto_generate_broll.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tarungupta/Making It Big/Claude/content-machine</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/auto_broll.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/auto_broll.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tarungupta/miniconda3/envs/content_engine_env/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/tarungupta</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
